# timers/batterylife.service
[Unit]
Description=Make notification when battery is below 10%

[Service]
Type=simple
ExecStart=/usr/bin/batterylife

